# binance-webhook-service/config/webhook_service.service
[Unit]
Description=Binance Futures Webhook Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/Algo-Trading/binance-webhook-service
Environment="PATH=/opt/Algo-Trading/binance-webhook-service/venv/bin"
EnvironmentFile=/opt/Algo-Trading/binance-webhook-service/.env
ExecStart=/opt/Algo-Trading/binance-webhook-service/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --threads 2 --timeout 30 --access-logfile /opt/Algo-Trading/binance-webhook-service/logs/access.log --error-logfile /opt/Algo-Trading/binance-webhook-service/logs/error.log --pythonpath /opt/Algo-Trading/binance-webhook-service/src binance_webhook_service:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

